# etc/supervisord/router.ini
[program:router]
directory=/srv/work/football-matome
command=cargo run --bin router
stdout_logfile_maxbytes=1MB
stderr_logfile_maxbytes=1MB
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.err
autorestart=true
stopsignal=INT
stopwaitsecs=60
user=ubuntu
numprocs=1
process_name=%(program_name)s_%(process_num)02d
